# ansible/roles/data-lifecycle/templates/failure-notify@.service.j2
# Copyright (c) 2025 MikeT LLC. All rights reserved.
#
# failure-notify@.service
# Generic failure notification service template
# Triggered via OnFailure= in other systemd units
#
# Usage: OnFailure=failure-notify@%p.service
# The %p placeholder captures the unit prefix (name without .service suffix)
#
# Deployed by Ansible data-lifecycle role

[Unit]
Description=Failure notification for %i
Documentation=https://github.com/miket-llc/miket-infra/docs/architecture/components/DATA_ESTATE_SPEC.md

[Service]
Type=oneshot

# Set environment for the script
Environment=DATA_LIFECYCLE_LOG_DIR={{ data_lifecycle_logs_dir | default('/space/_ops/logs/data-lifecycle') }}
Environment=DATA_LIFECYCLE_MARKERS_DIR={{ data_lifecycle_markers_dir | default('/space/_ops/data-estate/markers') }}

# Run the failure notification script with the failed unit name
ExecStart=/usr/local/bin/failure-notify.sh %i

User=root
StandardOutput=journal
StandardError=journal
SyslogIdentifier=failure-notify
