# ansible/roles/data-lifecycle/templates/flux-local.service.j2
[Unit]
Description=Local snapshot of /flux
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/flux-local-snap.sh
User=root

