# ansible/roles/data-lifecycle/templates/flux-graduate.service.j2
[Unit]
Description=Graduate cold data from /flux to /space
Wants=flux-graduate.timer

[Service]
Type=oneshot
ExecStart=/usr/local/bin/flux-graduate.sh
User=root

