# ansible/roles/data-lifecycle/templates/flux-backup.service.j2
[Unit]
Description=Backup /flux to Backblaze B2 (restic cloud)
Documentation=https://github.com/miket-llc/miket-infra/docs/architecture/components/DATA_ESTATE_SPEC.md
After=network-online.target
Wants=network-online.target
OnFailure={{ data_lifecycle_on_failure_target | default('failure-notify@%p.service') }}

[Service]
Type=oneshot

# Credentials from AKV via secrets_sync
EnvironmentFile={{ data_lifecycle_env_file }}

# Marker file directory for Data Estate collector
Environment=DATA_LIFECYCLE_MARKERS_DIR={{ data_lifecycle_markers_dir | default('/space/_ops/data-estate/markers') }}
Environment=DATA_LIFECYCLE_LOG_DIR={{ data_lifecycle_logs_dir | default('/space/_ops/logs/data-lifecycle') }}

# Per-host restic repo. restic tags snapshots by hostname so multiple
# hosts CAN safely share a repo, but we keep them separate for clarity
# and so retention policies can diverge. restic-check + restore-test
# pull from the same `data_lifecycle_restic_repo` var so all three
# point at the same destination.
Environment=RESTIC_REPO={{ data_lifecycle_restic_repo | default('b2:miket-backups-restic:flux') }}

ExecStart=/usr/local/bin/flux-backup.sh
User=root

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=flux-backup

# Timeout (generous for large backups)
TimeoutStartSec=7200

[Install]
WantedBy=multi-user.target

