# ansible/roles/data-lifecycle/templates/space-mirror.service.j2
[Unit]
Description=Mirror /space to Backblaze B2
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/miket/storage-credentials.env
ExecStart=/usr/local/bin/space-mirror.sh
User=root

