# ansible/roles/_deprecated/remote_server_linux_vnc/templates/tigervnc.service.j2
[Unit]
Description=TigerVNC x0vncserver (Display Sharing - GNOME Session)
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
User={{ ansible_user | default('mdt') }}
Group={{ ansible_user | default('mdt') }}
WorkingDirectory=/home/{{ ansible_user | default('mdt') }}
Environment="HOME=/home/{{ ansible_user | default('mdt') }}"
Environment="DISPLAY=:{{ vnc_display | default('0') }}"
Environment="XAUTHORITY={{ xauth_path | default('/run/user/1000/gdm/Xauthority') }}"
ExecStartPre=/bin/sh -c 'DISPLAY=:{{ vnc_display | default("0") }} XAUTHORITY={{ xauth_path | default("/run/user/1000/gdm/Xauthority") }} xset q &>/dev/null || sleep 2'
ExecStart=/usr/bin/x0vncserver -display :{{ vnc_display | default('0') }} -rfbport {{ remote_port | default(5900) }} -PasswordFile /home/{{ ansible_user | default('mdt') }}/.vnc/tigervnc-passwd -SecurityTypes VncAuth -AlwaysShared -localhost no -fg
Restart=on-failure
RestartSec=10
TimeoutStartSec=30

[Install]
WantedBy=default.target

