# ansible/roles/llm-serving/templates/vllm-refresh@.timer.j2
# {{ ansible_managed }}
[Unit]
Description=Nightly refresh of vLLM container (%i)

[Timer]
OnCalendar={{ vllm_refresh_schedule | default('daily') }}
RandomizedDelaySec=10m
Unit=vllm-refresh@%i.service

[Install]
WantedBy=timers.target
