# ansible/roles/ops_verify/templates/dmarc-summary.service.j2
[Unit]
Description=Daily DMARC aggregate-report summary for {{ ops_verify_email_domain }}
After=network-online.target
Wants=network-online.target
RequiresMountsFor=/space

[Service]
Type=oneshot
Environment=DMARC_DOMAIN={{ ops_verify_email_domain }}
Environment=NODE_EXPORTER_TEXTFILE_DIR={{ ops_verify_metrics_dir }}
Environment=OPS_VERIFY_LOG_DIR={{ ops_verify_log_dir }}
Environment=DMARC_REPORT_DIR=/space/_ops/reports/dmarc
ExecStart={{ ops_verify_bin_dir }}/dmarc-summary.sh
User=root

StandardOutput=journal
StandardError=journal
SyslogIdentifier=dmarc-summary

TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
