# ansible/roles/ops_verify/templates/secrets-drift.service.j2
[Unit]
Description=Secrets inventory drift probe (miket-infra .ops/secrets.yaml)
After=network-online.target
Wants=network-online.target
RequiresMountsFor=/space /flux

[Service]
Type=oneshot
# secrets_drift.py doesn't touch Azure — it's a pure YAML-to-severity
# evaluation, so the service runs as root (simplest ownership story for
# the shared log/metrics dirs). Tier 2 services that DO call `az` run as
# User={{ ops_verify_az_user }} so they can reach the user-level
# credential cache.
Environment=INFRA_DIR=/flux/ops/miket-infra
Environment=OPS_VERIFY_VENV=/flux/ops/ops-verify-venv
Environment=NODE_EXPORTER_TEXTFILE_DIR={{ ops_verify_metrics_dir }}
Environment=OPS_VERIFY_LOG_DIR={{ ops_verify_log_dir }}
Environment=SECRETS_DRIFT_REPORT_DIR=/space/_ops/reports/secrets
ExecStart={{ ops_verify_bin_dir }}/secrets-drift.sh

StandardOutput=journal
StandardError=journal
SyslogIdentifier=secrets-drift

TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
