# ansible/roles/ops_verify/templates/secrets-drift.timer.j2
[Unit]
Description=Run secrets-inventory drift probe every hour

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
RandomizedDelaySec=5min
Persistent=true

[Install]
WantedBy=timers.target
