# ansible/roles/nextcloud_server/templates/nextcloud-m365-sync.service.j2
# Copyright (c) 2025 MikeT LLC. All rights reserved.
# Managed by Ansible - do not edit directly

[Unit]
Description=Nextcloud M365 Ingestion Sync (OneDrive/SharePoint → /space)
Documentation=https://github.com/miket-llc/miket-infra/docs/architecture/components/DATA_ESTATE_SPEC.md
After=network-online.target podman.socket
Wants=network-online.target
Requires=podman.socket
OnFailure=failure-notify@%n.service

[Service]
Type=oneshot
ExecStart={{ nextcloud_bin_path }}/sync-m365.sh

# Credentials passed via script loading env file directly
# (safer handling of special characters in passwords)

# User context
User=root
Group=root

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nextcloud-m365-sync

# Timeout (allow up to 4 hours for large syncs)
TimeoutStartSec=14400

[Install]
WantedBy=multi-user.target

