# ansible/roles/nextcloud_server/templates/nextcloud-home-sweeper.service.j2
# Copyright (c) 2025 MikeT LLC. All rights reserved.
# Managed by Ansible - do not edit directly

[Unit]
Description=Nextcloud Home Sweeper - Detect stray files in internal user homes
Documentation=file:///home/mdt/miket-infra-devices/docs/guides/nextcloud_on_motoko.md
After=docker.service nextcloud.service

[Service]
Type=oneshot
ExecStart={{ nextcloud_bin_path }}/nextcloud-home-sweeper.sh
User=root
Group=root

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nextcloud-home-sweeper

