# ansible/roles/nextcloud_server/templates/nextcloud.service.j2
# Copyright (c) 2025 MikeT LLC. All rights reserved.
# Managed by Ansible - do not edit directly
#
# Nextcloud Stack Service (Podman)
# Includes post-start validation to handle podman DNS race conditions

[Unit]
Description=Nextcloud Stack (Podman Compose)
Requires=podman.socket
After=podman.socket network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory={{ nextcloud_workdir }}
EnvironmentFile={{ nextcloud_env_path }}

# Pre-start: ensure clean state and pull images
ExecStartPre=/usr/bin/podman compose pull

# Start the stack
ExecStart=/usr/bin/podman compose up -d --remove-orphans

# Post-start: wait for containers to be healthy, then validate DB connectivity
# This catches the podman DNS race condition and auto-recovers
ExecStartPost=/bin/bash -c 'sleep 15 && {{ nextcloud_bin_path }}/nextcloud-healthcheck.sh recover'

ExecStop=/usr/bin/podman compose down
ExecReload=/usr/bin/podman compose restart

# Restart policy - auto-recover from failures
Restart=on-failure
RestartSec=30

# Give enough time for health check and recovery
TimeoutStartSec=300

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nextcloud

[Install]
WantedBy=multi-user.target
