# ansible/roles/system-updates/templates/update-timer.service.j2
[Unit]
Description={{ update_timer_description }}
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User={{ ansible_user }}
WorkingDirectory={{ ansible_playbook_dir }}
# Note: No --ask-vault-pass for automated runs; update playbooks should not require vault
ExecStart=/usr/bin/ansible-playbook -i {{ update_inventory_path }} {{ update_playbook_path }}
StandardOutput=journal
StandardError=journal
# Allow up to 2 hours for updates across all hosts
TimeoutStartSec=7200

[Install]
WantedBy=multi-user.target

