# ansible/roles/claude_discord_agent/templates/claude-discord-agent.service.j2
# Claude Discord Agent — systemd user service
# Managed by Ansible (miket-infra-devices)
# Role: claude_discord_agent
#
# Uses a SEPARATE tmux socket (-L claude) so it never conflicts
# with your personal tmux sessions.
#
# Memory and bootstrap hook provide continuity across restarts.
# --resume is NOT used because it prevents --channels from initializing.
#
# Access:  tmux -L claude attach
# Detach:  Ctrl-b, d
# Logs:    journalctl --user -u claude-discord-agent -f

[Unit]
Description=Claude Code Discord Agent
Documentation=https://github.com/anthropics/claude-code
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=forking

# Brief delay so DNS/network are fully ready (Discord plugin needs connectivity)
ExecStartPre=/usr/bin/sleep {{ claude_agent_start_delay }}

ExecStart=/usr/bin/tmux -L claude new-session -d -s agent \
  {{ claude_agent_binary }} \
  --dangerously-skip-permissions \
  --channels {{ claude_agent_channel }}
ExecStop=/usr/bin/tmux -L claude kill-session -t agent
WorkingDirectory={{ claude_agent_workdir }}
Restart=always
RestartSec={{ claude_agent_restart_sec }}

# Environment — bun must be in PATH for the Discord plugin MCP server to start
Environment=HOME={{ claude_agent_home }}
Environment=PATH={{ claude_agent_home }}/.bun/bin:{{ claude_agent_home }}/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=TERM=xterm-256color
{% for key, value in claude_agent_extra_env.items() %}
Environment={{ key }}={{ value }}
{% endfor %}

[Install]
WantedBy=default.target
